# openbsd_rc.d_pc_fcgi
#!/bin/sh

daemon="/home/kyle/src/pricechart/pc_fcgi"

. /etc/rc.d/rc.subr

rc_cmd $1
