# etc/openbsd_rc.d_ps_fcgi
#!/bin/sh

daemon="perl -T /home/kyle/src/pricesloth/ps_fcgi"

. /etc/rc.d/rc.subr

pexp="perl: ps_fcgi sloth"

rc_cmd $1
